# config/runit/uwsgi/run
#!/bin/bash -e
exec uwsgi /app/config/uwsgi.ini
